# hello-crontab-project/crontab
#run python script every minutes
*/2 * * * * python /app/app.py >> /logfile